# docker/deploy.ini
[program:opsany-saas-deploy]
command=gunicorn wsgi --bind 0.0.0.0:7003 -k gevent -w 4 -n deploy --access-logfile - --error-logfile -
directory=/opt/opsany/deploy/
environment=DJANGO_SETTINGS_MODULE=settings,BK_ENV="production",PAAS_LOGGING_DIR="/opt/opsany/logs"
startsecs=0
stopwaitsecs=0
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/opt/opsany/logs/saas_deploy.log
